# frontend-vm/custom-nginx.service
[Unit]
Description=Custom Nginx Docker Container
After=docker.service
Requires=docker.service

[Service]
Restart=always
RestartSec=5

ExecStartPre=/usr/bin/docker rm -f custom-nginx || true
ExecStart=/usr/bin/docker run \
    --name custom-nginx \
    -p 80:80 \
    custom-nginx:1.0

ExecStop=/usr/bin/docker stop custom-nginx

[Install]
WantedBy=multi-user.target
